# external-monitor/95-monitor-switch.rules
ACTION=="change", KERNEL=="card1", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/mlp/.Xauthority", RUN+="/home/mlp/bin/autorandr.sh"
